# Services/user/restart_showcontrol.service
[Unit]
Description=Restart showcontrol service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart showcontrol.service
