# Services/sc_mix.service.in
[Unit]
Description=SuperCollider mixer and encoder
After=sound.target local-fs.target seamless_jack.service

[Install]
WantedBy=default.target

[Service]
Type=idle
PrivateTmp=true
Environment="QT_QPA_PLATFORM=offscreen"
ExecStart=/usr/local/bin/sclang @share_dir@/SC/HuFo_SERVER.scd
LimitRTPRIO=95
LimitRTTIME=infinity
LimitMEMLOCK=infinity
